# crontab
LOGDIR=/home/oracle/scripts/logs
DBNAME=ORADB
1 20 * * * /home/oracle/scripts/rmanctl.sh $DBNAME  > "$LOGDIR"/$DBNAME-$(date +\%d-\%m-\%Y).out
